# systemd/atomic-arch-scan.timer
[Unit]
Description=Run Atomic Arch AUR scan weekly

[Timer]
OnCalendar=weekly
Persistent=true
# Spread load across hosts that share the same weekly schedule.
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
