# systemd/atomic-arch-notify@.service
[Unit]
Description=Atomic Arch AUR scan with notification hook (%i)
After=network-online.target

[Service]
Type=oneshot
# Template unit (%i) — use instances like atomic-arch-notify@hostname for per-host hooks.
Environment=AUR_RESPONSE_DIR=/opt/atomic-arch-response-toolkit
ExecStart=/usr/bin/fish %AUR_RESPONSE_DIR%/run.fish --local --report --json
# Replace ExecStartPost with mail, webhook, or notify-send; $? is the scan exit code.
ExecStartPost=/usr/bin/systemd-cat echo "Atomic Arch scan unit %i finished with exit code $?"

[Install]
WantedBy=multi-user.target
