# systemd/aur-response-notify@.service
[Unit]
Description=AUR incident scan with notification hook (%i)
After=network-online.target

[Service]
Type=oneshot
# Template unit (%i) — use instances like aur-response-notify@hostname for per-host hooks.
Environment=AUR_RESPONSE_DIR=/opt/aur-response-toolkit
ExecStart=/usr/bin/fish %AUR_RESPONSE_DIR%/run.fish --local --report --json
# Replace ExecStartPost with mail, webhook, or notify-send; $? is the scan exit code.
ExecStartPost=/usr/bin/systemd-cat echo "AUR incident scan unit %i finished with exit code $?"

[Install]
WantedBy=multi-user.target
